# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:frontend]
command=npm run build 
command=serve -s build -l 3000
directory=/app/website
autostart=true
autorestart=true
redirect_stderr=true

[program:django]
command=python3 manage.py makemigrations
command=python3 manage.py migrate
command=python3 manage.py runserver   0.0.0.0:8000
directory=/app/djangoapp/backend
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true

